# shairport-sync-src/rootfs/etc/services.d/avahi/run
#!/usr/bin/with-contenv bashio

declare -a options

bashio::log.info 'Starting Avahi...'

# wait for dbus
until [[ -e /var/run/dbus/system_bus_socket ]];
do
	sleep 1
done

exec avahi-daemon