# shairport-sync-src/rootfs/etc/services.d/shairport-sync/run
#!/usr/bin/with-contenv bashio
declare -a options

bashio::log.info 'Starting the Shairport Sync daemon...'

# Debug mode handler
#if bashio::debug; then
#  options+=(--verbose)
#fi

# Run shairport-sync
exec shairport-sync "${options[@]}"
