# shairport-sync-src/rootfs/etc/services.d/dbus/run
#!/usr/bin/with-contenv bashio

declare -a options

bashio::log.info 'Starting the D-Bus daemon...'
exec dbus-daemon --system --nofork